# ops/systemd/kitzchat.service
[Unit]
Description=KitzChat (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Template default. Prefer deploying the repo to /opt/kitzchat (or similar)
# instead of user home directories.
WorkingDirectory=/opt/kitzchat

# Secrets belong in an EnvironmentFile, not in the unit file.
# Point this at your deployment-specific env file.
EnvironmentFile=-/etc/kitzchat/kitzchat.env

# Optional: path to a 1Password op:// env template file (non-secret).
Environment=KITZCHAT_OP_ENV_FILE=/etc/kitzchat/kitzchat.op.env
# Optional: 1Password runtime mode (off|auto|required). Default is auto.
Environment=KITZCHAT_1PASSWORD_MODE=auto

ExecStart=/opt/kitzchat/scripts/start-standalone.sh
Restart=on-failure
RestartSec=3

# Run as a dedicated user in template deployments.
User=kitzchat
Group=kitzchat

# Hardening baseline. Adjust ReadWritePaths to match your KITZCHAT_DB_PATH/KITZCHAT_STATE_DIR.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/kitzchat /var/lib/kitzchat

[Install]
WantedBy=multi-user.target
